# tmp/droid-config.files
%defattr(-,root,root,-)
/usr/lib/sysctl.d/ipv6-settings.conf
/usr/lib/sysctl.d/ipv4-settings.conf
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/tmpfiles.d/droid-hal.conf
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-load-firmware.sh
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-init-done.sh
/etc/lowi.conf
/etc/selinux/config
/etc/media_codecs_performance.xml
/etc/xdg/QtProject/QtQuick2.conf
/etc/xdg/QtProject/Sensors.conf
/etc/xdg/QtProject/QPlatformTheme.conf
/etc/media_codecs_google_video.xml
/etc/izat.conf
/etc/mixer_paths.xml
/etc/sec_config
/etc/media_codecs_google_telephony.xml
/etc/media_codecs.xml
/etc/xtra_root_cert.pem
/etc/obexd/opp.conf
/etc/obexd/root
/etc/obexd/pbap.conf
/etc/obexd/excludes/bluetooth:pcsuite
/etc/obexd/noplugins/ftp
/etc/obexd/noplugins/irmc
/etc/obexd/noplugins/mas
/etc/obexd/noplugins/syncevolution
/etc/sensorfw/primaryuse.conf
/etc/modprobe.d/blacklist-evbug.conf
/etc/ofono/binder.conf
/etc/ofono/binder.d/dual-sim.conf
/etc/media_profiles.xml
/etc/media_codecs_google_audio.xml
/etc/gps.conf
/etc/usb-moded/dyn-modes/mtp_mode-droid.ini
/etc/xtwifi.conf
/etc/audio_platform_info.xml
/etc/sap.conf
/etc/media_codecs_ffmpeg.xml
/etc/acdbdata
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/nemo/60-multimedia.conf
/var/lib/environment/ofono/noplugin.conf
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/rules.d/998-droid-system.rules
/usr/share/csd/settings.d/hw-settings.ini
/usr/share/icons/hicolor/86x86/apps/harbour-sailfishpen.png
/usr/share/harbour-spen-menu/qml/cover/CoverPage.qml
/usr/share/harbour-spen-menu/qml/pages/Main.qml
/usr/share/harbour-spen-menu/qml/pages/About.qml
/usr/share/harbour-spen-menu/qml/harbour-spen-menu.qml
/usr/share/harbour-spen-menu/translations/harbour-spen-menu.ts
/usr/share/harbour-spen-menu/translations/harbour-spen-menu-de.ts
/usr/libexec/droid-hybris/system/etc/init/racoon.rc
/usr/libexec/droid-hybris/system/etc/init/gatekeeperd.rc
/usr/libexec/droid-hybris/system/etc/init/vold.rc
/usr/libexec/droid-hybris/system/etc/init/dumpstate.rc
/usr/libexec/droid-hybris/system/etc/init/bootstat.rc
/usr/libexec/droid-hybris/system/etc/init/logd.rc
/usr/libexec/droid-hybris/system/etc/init/drmserver.rc
/usr/libexec/droid-hybris/system/etc/init/mediaextractor.rc
/usr/libexec/droid-hybris/system/etc/init/vdc.rc
/usr/libexec/droid-hybris/system/etc/init/mediadrmserver.rc
/usr/libexec/droid-hybris/system/etc/init/perfprofd.rc
/usr/libexec/droid-hybris/system/etc/init/rild.rc
/usr/libexec/droid-hybris/system/etc/init/superuser.rc
/usr/libexec/droid-hybris/system/etc/init/logcatd.rc
/usr/libexec/droid-hybris/system/etc/init/uncrypt.rc
/usr/libexec/droid-hybris/system/etc/init/debuggerd.rc
/usr/libexec/droid-hybris/system/etc/init/mediacodec.rc
/usr/libexec/droid-hybris/system/etc/init/atrace.rc
/usr/libexec/droid-hybris/system/etc/init/init-debug.rc
/usr/libexec/droid-hybris/system/etc/init/mtpd.rc
/usr/libexec/droid-hybris/system/etc/init/installd.rc
/usr/lib/mozembedlite/chrome/embedlite/content/duckduckgo.xml
/usr/lib/mozembedlite/chrome/embedlite/content/qwant.xml
/usr/lib/mozembedlite/chrome/embedlite/content/hulbee.xml
/usr/lib/mozembedlite/chrome/embedlite/content/startpage.xml
/usr/lib/mozembedlite/chrome/embedlite/content/disroot.xml
/usr/lib/mozembedlite/chrome/embedlite/content/seznam-cz.xml
/usr/lib/mozembedlite/chrome/embedlite/content/baidu.xml
/usr/lib/mozembedlite/chrome/embedlite/content/searx-https.xml
/usr/lib/mozembedlite/chrome/embedlite/content/ecosia.xml
/usr/lib/mozembedlite/chrome/embedlite/content/ixquick.xml
/usr/lib/mozembedlite/chrome/embedlite/content/swisscows.xml
/usr/lib/systemd/system/spen.service
/usr/lib/systemd/system/actdead-pre.target.wants/droid-hal-init.service
/usr/bin/droid/wifi-tidy.sh
/usr/bin/droid/evtest
/usr/bin/droid/wifi-save-status.sh
/usr/bin/droid/spen-start-service.sh
/usr/bin/droid/lxc-scripts/debian.sh
/usr/bin/droid/lxc-scripts/debian-mouse.sh
/usr/bin/droid/spen-daemon.sh
/usr/bin/droid/wifi-on.sh
/usr/bin/droid/wifi.conf
/usr/bin/logcat
/usr/bin/harbour-spen-notes
/usr/bin/harbour-spen-menu
/usr/bin/clear
/usr/bin/hotspot.sh
/usr/bin/kimmoli-diyicons.sh
/usr/local/share/info/find.info-2
/usr/local/share/info/find-maint.info
/usr/local/share/info/find.info-1
/usr/local/share/info/find.info
/usr/local/share/locale/be/LC_MESSAGES/findutils.mo
/usr/local/share/man/man1/updatedb.1
/usr/local/share/man/man1/find.1
/usr/local/share/man/man1/xargs.1
/usr/local/share/man/man1/locate.1
/usr/local/share/man/man5/locatedb.5
/usr/local/libexec/frcode
/usr/local/var/locatedb
/usr/local/bin/find
/usr/local/bin/locate
/usr/local/bin/xargs
/usr/local/bin/updatedb
/etc/mce/99evdev.ini
/etc/mce/60-battery-udev-devices.ini
/etc/mce/60-rgb.ini
/etc/ofono/ril_subscription.conf
/etc/systemd/system/multi-user.target.wants/spen.service
/etc/firmware/bcm4358a1.conf
/etc/sysconfig/bluetooth-rfkill-event
/etc/sysconfig/connman
/var/lib/environment/compositor/droid-hal-device.conf
/var/lib/environment/ofono/noplugin.conf
/lib/udev/rules.d/100-wifi-on.rules
/lib/udev/rules.d/101-wifi-off.rules
/usr/share/ssu/features.d/adaptation-community.ini
