# sparse/lib/udev/rules.d/100-wifi-on.rules
# do not edit this file, it will be overwritten on update
# edp17 - save wifi status when wifi is toggle is used (rfkill changes status)

ACTION=="change", SUBSYSTEM=="rfkill", ATTR{type}=="wlan", ATTR{state}=="1", RUN+="/usr/bin/droid/wifi-save-status.sh on"
