# sparse/lib/udev/rules.d/101-discharging.rules
# edp17 - Detect when cable removed

ACTION=="change", SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", RUN+="/usr/bin/droid/discharging.sh"
