# sparse/lib/udev/rules.d/103-bluetooth-brcm-on.rules
# do not edit this file, it will be overwritten on update
# edp17 - start brcm_patchram_plus when rfkill changes status

ACTION=="change", SUBSYSTEM=="rfkill", ATTR{type}=="bluetooth", ATTR{state}=="1", RUN+="/usr/bin/droid/bt-start.sh"
