# sparse/usr/lib/systemd/system/spen.service
[Unit]
Description=S-Pen listener

[Service]
ExecStart=/usr/bin/droid/spen-start-service.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
